# systemd/doublethink.service
[Unit]
Description=doublethink secure pub/sub broker
Documentation=https://ra-yavuz.github.io/doublethink/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The channels endpoint is public; the admin/pairing endpoint binds to loopback
# only and must never be exposed off-host (it is the unauthenticated control
# surface, by the "whoever reaches localhost controls this host" model).
ExecStart=/usr/bin/doublethink serve --addr :8080 --admin-addr 127.0.0.1:8081 --state /var/lib/doublethink/state.json
Restart=on-failure
RestartSec=2

# Hardening: doublethink needs no privilege beyond binding its port and writing
# its state file. Lock the rest down.
DynamicUser=yes
StateDirectory=doublethink
RuntimeDirectory=doublethink
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
