# debian/postrm
#!/bin/sh
# postrm for doublethink
set -e

case "$1" in
    remove|disappear)
        if [ -d /run/systemd/system ]; then
            systemctl stop doublethink.service 2>/dev/null || true
            systemctl disable doublethink.service 2>/dev/null || true
            systemctl daemon-reload || true
        fi
        ;;
    purge)
        # StateDirectory under DynamicUser lives at /var/lib/doublethink; remove
        # the channel registry on purge. This deletes the authorised-key state, so
        # paired channels must be re-created after a purge+reinstall (deliberate).
        rm -rf /var/lib/doublethink
        ;;
    upgrade|failed-upgrade|abort-install|abort-upgrade)
        ;;
    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
